# configuration/myturn@.service
[Unit]
Description=%p, %i daemon
After=network.target

[Service]
WorkingDirectory=/var/www/%p-%i
Environment=SERVER_PORT=3000
ExecStart=/usr/bin/node api.js --local
Restart=on-failure
SyslogIdentifier=%p-%i

[Install]
WantedBy=multi-user.target
